# sub-agent.service
[Unit]
Description=OpenClaw Sub-Agent Telegram Interface
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=boo
WorkingDirectory=/home/boo/sub-agent-interface
Environment=CLAW_TELEGRAM_TOKEN=YOUR_TOKEN_HERE
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/boo/sub-agent-interface/sub-agent-interface.py --serve
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claw-sub-agent

[Install]
WantedBy=multi-user.target
